# roles/setup_system/templates/pi-visual.ini.j2
## GENERATED BY ANSIBLE

[program:pi-visual]
# refer to: http://supervisord.org/configuration.html#program-x-section-settings
user=pnguyen
directory={{ repo_path }}
command={{ repo_path }}/.venv/bin/gunicorn app:app -c {{ repo_path }}/gunicorn.py
autorestart=true
autorestart=true        ; runs on startup
startretries=3
logfile=/dev/null
